# scripts/gymnasticon-ant.service
[Unit]
Description=Gymnasticon ANT+ Bridge
After=network.target

[Service]
Type=simple
# nvm is a shell function; use absolute path to node 12
ExecStart=/home/jeffreysnowden/.nvm/versions/node/v12.22.12/bin/node /home/jeffreysnowden/gymnasticon/scripts/ant-bridge.js
Environment=NODE_PATH=/home/jeffreysnowden/gymnasticon/node_modules
WorkingDirectory=/home/jeffreysnowden/gymnasticon
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
