# scripts/gymnasticon-ble.service
[Unit]
Description=Gymnasticon BLE Reader
After=bluetooth.target gymnasticon-ant.service
Requires=bluetooth.service
Wants=gymnasticon-ant.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPre=/bin/sleep 3
ExecStart=/home/jeffreysnowden/bleak-env/bin/python /home/jeffreysnowden/gymnasticon/scripts/ble-reader.py
WorkingDirectory=/home/jeffreysnowden/gymnasticon
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
